# scripts/systemd/celery-worker.service
[Unit]
Description=Zen-AI-Pentest Celery Worker
After=network.target redis.service postgresql.service

[Service]
Type=simple
User=zenpentest
Group=zenpentest
WorkingDirectory=/opt/zen-ai-pentest
Environment=PYTHONPATH=/opt/zen-ai-pentest
Environment=REDIS_URL=redis://localhost:6379/0
Environment=DATABASE_URL=postgresql://postgres:password@localhost:5432/zen_pentest
ExecStart=/opt/zen-ai-pentest/venv/bin/celery -A celery_app worker -l info -c 4 -Q celery,scans.critical,scans.high,scans.normal,scans.low,notifications,reports
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
