# scripts/systemd/celery-beat.service
[Unit]
Description=Zen-AI-Pentest Celery Beat Scheduler
After=network.target redis.service postgresql.service celery-worker.service

[Service]
Type=simple
User=zenpentest
Group=zenpentest
WorkingDirectory=/opt/zen-ai-pentest
Environment=PYTHONPATH=/opt/zen-ai-pentest
Environment=REDIS_URL=redis://localhost:6379/0
Environment=DATABASE_URL=postgresql://postgres:password@localhost:5432/zen_pentest
ExecStart=/opt/zen-ai-pentest/venv/bin/celery -A celery_app beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
